# server/com.wheel.image-server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  launchd agent for the Wheel-of-the-Year ecology image service.
  Install on the Mac mini:
    cp server/com.wheel.image-server.plist ~/Library/LaunchAgents/
    launchctl load ~/Library/LaunchAgents/com.wheel.image-server.plist
  Reload after changes (re-reads the plist):
    launchctl unload ~/Library/LaunchAgents/com.wheel.image-server.plist
    launchctl load   ~/Library/LaunchAgents/com.wheel.image-server.plist
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.wheel.image-server</string>

  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/node</string>
    <string>/Users/bradfordadams/Sites/wheel-of-the-year/server/image-server.mjs</string>
  </array>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PORT</key>
    <string>7871</string>
    <key>IMAGE_CACHE_DIR</key>
    <string>/Users/bradfordadams/Sites/wheel-of-the-year/image-cache</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/bradfordadams/Library/Logs/wheel-image-server.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/bradfordadams/Library/Logs/wheel-image-server.log</string>
</dict>
</plist>
